# packages/flink/flink.service
[Unit]
Description=Apache Flink 1.9.3
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/flink-1.9.3/bin/start-cluster-on-8082.sh
ExecStop=/usr/sbin/flink-1.9.3/bin/stop_cluster.sh

[Install]
WantedBy=multi-user.target
