# packages/log4shell/log4shell.service
[Unit]
Description=Log4Shell vulnerable server

[Service]
ExecStart=/usr/bin/java -Dcom.sun.jndi.ldap.object.trustURLCodebase=true -jar /usr/sbin/log4shell.jar
Type=simple
WorkingDirectory=/usr/sbin/

User=e1-target
Group=e1-target

[Install]
WantedBy=multi-user.target
